# scripts/unifi/systemd/unifi-export.timer
[Unit]
Description=Run UniFi Data Export daily

[Timer]
# Run daily at 03:30 with a small randomized delay
OnCalendar=*-*-* 03:30:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
